# etc/init.d/hermes
#!/sbin/openrc-run
name="hermes"
description="Hermes AI agent (youself.io)"
command="python3"
command_args="-m hermes --transport youself_gateway"
command_user="root"
command_background=true
pidfile="/run/hermes.pid"
command_stdout="/var/log/hermes/agent.log"
command_stderr="/var/log/hermes/agent.err"

start_pre() {
    [ -f /etc/openclaw/env ] && export $(grep -v '^#' /etc/openclaw/env | xargs)
    export OPENAI_API_BASE="${LLM_PROXY_BASE_URL:-}"
    export OPENAI_API_KEY="${LLM_PROXY_API_KEY:-}"
}

depend() {
    after net
}
